# VoidLinux/etc/sv/etserver/run
#!/bin/sh
exec 2>&1
[ -r conf ] && . ./conf

exec /home/linuxbrew/.linuxbrew/opt/et/bin/etserver $OPTS

